# etcd/scripts/deb/after-install.sh
# if [ $1 -eq 1 ] ; then
        # Initial installation

        update-rc.d etcd defaults >/dev/null 2>&1 || :
# fi
